# scripts/celery-worker.conf
; ==================================
;  celery worker supervisor
; ==================================

[program:flapisk]
directory=/var/www/flapisk/backend
command=/var/www/flapisk/backend/env/bin/celery worker -A application.celery --loglevel=INFO

user=root
numprocs=1
stdout_logfile=/var/www/flapisk/logs/celery/access.log
stderr_logfile=/var/www/flapisk/logs/celery/error.log
stdout_logfile_maxbytes=5000000
stderr_logfile_maxbytes=5000000
stdout_logfile_backups=5000000
stderr_logfile_backups=5000000
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 36000

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000
